# net-stabilization.service
[Unit]
Description=Net Stabilization - Mining Fleet Power Control
Documentation=https://github.com/dkk/net-stabilization
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/dkk/Documents/GitHub/net-stabilization
ExecStartPre=/usr/bin/docker compose pull --quiet --ignore-pull-failures
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
ExecReload=/usr/bin/docker compose restart
TimeoutStartSec=300
TimeoutStopSec=60

# Restart policy
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
